# files/system/common/etc/profile.d/neofetch.sh
alias neofetch='/usr/bin/fastfetch -c /usr/share/ublue-os/corse/fastfetch.jsonc'
alias neowofetch='/usr/bin/fastfetch -c /usr/share/ublue-os/corse/fastfetch.jsonc'
alias hyfetch='/usr/bin/fastfetch -c /usr/share/ublue-os/corse/fastfetch.jsonc'
alias fastfetch='/usr/bin/fastfetch -c /usr/share/ublue-os/corse/fastfetch.jsonc'